# scripts/tpv-hardware.service
[Unit]
Description=TPV Hardware Printer Bridge
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/carbones_tpv
EnvironmentFile=/opt/carbones_tpv/.env
ExecStart=/opt/carbones_tpv/venv/bin/python /opt/carbones_tpv/scripts/hardware_bridge.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
